# services.d/nginx/run
#!/usr/bin/with-contenv bash

exec deepspeech-server --config /data/config.json
